# trunk/wxCron/crontab
#M    S   T M W    Befehl
* */7 * * *    ..\blackfiskDd.exe D:\Garage\projekte\blackfisk\trunk\MyOrg.bfp